# deploy/server/infra/systemd/itmen-pipeline-kaiten-poll.service
[Unit]
Description=ITMen Pipeline Kaiten inbound poll
After=network.target pb-itmen-pipeline.service itmen-pipeline-api.service

[Service]
Type=oneshot
WorkingDirectory=/opt/itmen-pipeline
EnvironmentFile=/opt/itmen-pipeline/.env
ExecStart=/usr/bin/node /opt/itmen-pipeline/scripts/run-kaiten-poll.js
Nice=10
